# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:airflow]
command=/usr/bin/bash -c "/root/run_airflow.sh"